# crontab
* * * * * ./env/bin/python recorder.py >>log.recorder 2>&1
