# crontab
* * * * * BLINKA_MCP2221=1 ./env/bin/python recorder.py >/dev/null
