# crontab
* * * * * ./env/bin/python recorder.py >/dev/null
